# srv.service
[Unit]
Description=ProdCal (Production Calendar + Book Production Pipeline)
After=network.target

[Service]
# Type=notify pairs with the binary's daemon.SdNotify(SdNotifyReady) call
# in srv.Server.Serve(), which fires after the TCP listener is bound.
# This gives systemd a positive readiness signal and prevents the MainPID
# race we saw under Type=simple where systemd could spawn a second start
# attempt before recognising the first one had succeeded.
Type=notify
NotifyAccess=main
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/prodcal
ExecStart=/home/exedev/prodcal/prodcal -listen :8000
# Restart on failure (non-zero exit) only — Restart=always combined with
# main.go's previous error-swallowing exit-0 path was the trigger for the
# orphan storm we hit during cutover. main.go now exits 1 on errors.
Restart=on-failure
RestartSec=5
Environment=HOME=/home/exedev
Environment=USER=exedev
# typesettingRoot() in srv/books.go resolves typesetting/ relative to
# WorkingDirectory by default; set JDBBS_TYPESETTING_DIR here if the
# binary is ever run from a different location.
Environment=JDBBS_TYPESETTING_DIR=/home/exedev/prodcal/typesetting
EnvironmentFile=-/home/exedev/prodcal/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
